# tests/services/rootfs/etc/ignity/services/no-start-at-runtime/run
#!/usr/bin/with-env bash
echo "no-start service at runtime"
sleep infinity
