# tests/services/rootfs/etc/ignity/services/default/run
#!/usr/bin/with-env bash
echo "default service"
sleep infinity
