# src/etc/s6/services/.s6-svscan/finish
#!/bin/execlineb -S0
/etc/s6/boot/stage3 $@