# tests/services/rootfs/etc/ignity/services/no-start-at-build-time/run
#!/usr/bin/with-env bash
echo "no-start service at build time"
sleep infinity
